# crontab
# Run the script every 7 days
0 0 */7 * * cd /usr/src/app && node index.js >> /var/log/cron.log 2>&1
